# scripts/cctv-watchdog.service
[Unit]
Description=cctv-viewer freeze watchdog (one-shot)
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/freeze-watchdog.sh
