# scripts/cctv-watchdog.timer
[Unit]
Description=Run the cctv-viewer freeze watchdog every 5 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
AccuracySec=30s
Unit=cctv-watchdog.service

[Install]
WantedBy=timers.target
